# cron/crontab
50 14 * * * root /usr/bin/python /app/cron/notify_debts.py >> /var/log/cron.log 2>&1

